# dz-reward.timer
[Unit]
Description=Run DZ-Reward daily

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
